# mullosc+mongodb/supervisord.conf
[supervisord]
nodaemon=true

[program:mongodb]
command=/usr/bin/mongod -f /etc/mongod.conf

[program:mollusc]
directory=/mollusc
command=python manage.py runserver 0.0.0.0:8080
